# crontab
0 4 * * * www-data /usr/local/bin/php /var/www/html/artisan freetrealcommand >> /dev/stdout 2>&1
0 * * * * www-data /usr/local/bin/php /var/www/html/artisan send-post-stay-emails >> /dev/stdout 2>&1